# etc/init.d/trimd
#!/sbin/openrc-run

CONFIG='_CFGPATH_'

checkconfig()
{
    if [ ! -f ${CONFIG} ] ; then
        eerror "${CONFIG} does not exists."
        return 1
    fi

    if [ -z "${DAEMON}" ] \
        || [ -z "${PID}" ] \
        || [ -z "${NICE}" ] \
        || [ -z "${IONICE}" ]
    then
        eerror "You must set config options in ${CONFIG} first"
        return 1
    fi

    if [ ! -f ${DAEMON} ] ; then
        eerror "Script ${DAEMON} does not exists."
        return 1
    fi

    if [ -z "${LOG}" ]; then
        LOGOPT=( --stdout-logger '/usr/bin/logger --tag trimd' )
    else
        LOGOPT=( --stdout ${LOG} )
    fi
}

start()
{
    ebegin "Starting trimd"
    checkconfig || return 1
    start-stop-daemon --start \
      --exec ${DAEMON} \
      --pidfile ${PID} \
      --background \
      --make-pidfile \
      --ionice ${IONICE} \
      --nice ${NICE} \
      "${LOGOPT[@]}" \
      -- -c ${CONFIG}
    eend $?
}

stop()
{
    ebegin "Stopping trimd"
    start-stop-daemon --stop \
      --exec ${DAEMON} \
      --pidfile ${PID} \
      -- -c ${CONFIG}
    eend $retval
}
